# wpm/etc/init.d/mariadb.ini
[program:mariadb]

command = /usr/bin/mysqld --user=mysql --datadir=/var/lib/mysql --log-error=/var/log/mysql/error.log

numprocs = 1
autostart = true
autorestart = true
