# wpm/etc/init.d/nginx.ini
[program:nginx]
command = /usr/sbin/nginx -c /home/wordpress/conf.d/nginx.conf
autostart = true
autorestart = true
