# wpm/etc/init.d/php-fpm.ini
[program:php-fpm]
command = /usr/bin/php-fpm -c /etc/php/php.ini -y /home/example.com/conf.d/php-fpm.conf
autostart = true
autorestart = true

