# wpm/etc/run/php-fpm.ini
[program:php-fpm]

command = /usr/bin/php-fpm -c /etc/php/php.ini -y /etc/wpm/php-fpm.conf

autostart = true
autorestart = true
